# roles/homeassistant/templates/homeassistant.service.j2
# This is a simple service file for systems with systemd to tun HA as user.
#
# For details please check https://home-assistant.io/getting-started/autostart/
#
[Unit]
Description=Home Assistant for %i
After=network.target

[Service]
Type=simple
User=%i
# Enable the following line if you get network-related HA errors during boot
ExecStartPre=/bin/sleep 60
# Use `whereis hass` to determine the path of hass
ExecStart=/usr/local/bin/hass --runner -c {{ hass_config_dir }}
SendSIGKILL=no
RestartForceExitStatus=100

[Install]
WantedBy=multi-user.target
